# capture/com.atriveo.cortex-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- screenpipe — AI that knows everything you've seen, said, or heard -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.atriveo.cortex-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__CAPTURE_DIR__/run-cortex-sync.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__CAPTURE_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>USER</key>
    <string>__USER__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>CORTEX_DRIVE_ROOT</key>
    <string>/Volumes/Kasliwal v2</string>
    <key>CORTEX_REPO_ROOT</key>
    <string>__REPO_ROOT__</string>
  </dict>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/Atriveo/cortex-sync.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/Atriveo/cortex-sync.launchd.log</string>
</dict>
</plist>
