# cloudflare/com.atriveo.cortex.cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.atriveo.cortex.cloudflared</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--config</string>
    <string>/Volumes/Kasliwal v2/working-memory/cloudflare/tunnel-host.yml</string>
    <string>run</string>
    <string>cortex</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/cortex-cloudflared.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/cortex-cloudflared.err</string>
</dict>
</plist>
